# ops/systemd/mascloner-api.service
[Unit]
Description=MasCloner API Server
Documentation=https://github.com/ali-raaed31/mascloner
After=network.target
Wants=network.target

[Service]
Type=exec
User=mascloner
Group=mascloner
WorkingDirectory=/srv/mascloner
Environment=PATH=/srv/mascloner/.venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/srv/mascloner
EnvironmentFile=/srv/mascloner/.env
ExecStart=/srv/mascloner/.venv/bin/python -m app.api.main
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
KillMode=mixed
TimeoutStopSec=30

# Security settings
NoNewPrivileges=false
PrivateTmp=yes
ProtectSystem=no
ProtectHome=yes
ReadWritePaths=/srv/mascloner/data /srv/mascloner/logs /srv/mascloner/etc
ProtectKernelTunables=no
ProtectKernelModules=no
ProtectControlGroups=no
RestrictRealtime=no
RestrictSUIDSGID=no
LockPersonality=no
MemoryDenyWriteExecute=no

# Capabilities - relaxed for rclone operations
CapabilityBoundingSet=CAP_SETUID CAP_SETGID
AmbientCapabilities=CAP_SETUID CAP_SETGID
PrivateDevices=no
ProtectKernelLogs=no
ProtectClock=no

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mascloner-api

[Install]
WantedBy=multi-user.target
